# packaging/weatherdress.service.in
[Unit]
Description=WeatherDress — affichage météo habillé
After=network.target graphical.target

[Service]
# SDL / Pygame : XDG_RUNTIME_DIR doit exister (souvent absent si pas de session graphique au boot).
ExecStartPre=/usr/bin/install -d -m 0700 -o @RUN_USER@ -g @RUN_USER@ @XDG_RUNTIME_DIR@
ExecStart=/usr/bin/python3 -m weatherdress.main
WorkingDirectory=@INSTALL_DIR@
Environment=PYTHONPATH=@INSTALL_DIR@/src
Restart=on-failure
RestartSec=10
User=@RUN_USER@
Environment=DISPLAY=:0
Environment=XAUTHORITY=@RUN_HOME@/.Xauthority
Environment=XDG_RUNTIME_DIR=@XDG_RUNTIME_DIR@
Environment=XDG_SESSION_TYPE=x11

[Install]
WantedBy=multi-user.target
