# packaging/weatherdress.service.in
[Unit]
Description=WeatherDress — affichage météo habillé
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m weatherdress.main
WorkingDirectory=@INSTALL_DIR@
Environment=PYTHONPATH=@INSTALL_DIR@/src
Restart=on-failure
RestartSec=10
User=@RUN_USER@
Environment=DISPLAY=:0
Environment=XAUTHORITY=@RUN_HOME@/.Xauthority
Environment=XDG_RUNTIME_DIR=@XDG_RUNTIME_DIR@
Environment=XDG_SESSION_TYPE=x11

[Install]
WantedBy=multi-user.target
