# packaging/weatherdress.service
[Unit]
Description=WeatherDress — affichage météo habillé
After=network.target

[Service]
ExecStart=/usr/bin/python3 -m weatherdress.main
WorkingDirectory=/home/pi/weatherdress
Environment=PYTHONPATH=/home/pi/weatherdress/src
Restart=on-failure
RestartSec=10
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

[Install]
WantedBy=multi-user.target
